# app/s6/teleport/run
#!/command/execlineb -P
exec teleport start -c /etc/teleport/teleport.yaml