# app/s6/notify/run
#!/bin/bash
inotifywait -m -r -q --event modify /etc/teleport/teleport.yaml | while read FILE
do
  echo "something happened on path $FILE"
  pkill -f teleport
done